# mediary.service
# ─────────────────────────────────────────────────────────────────────────────
# meDiary systemd service (user service)
# Installieren:
#   cp mediary.service ~/.config/systemd/user/mediary.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now mediary
#
# Logs:  journalctl --user-unit=mediary -f
# Stop:  systemctl --user stop mediary
# Restart: systemctl --user restart mediary
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=meDiary HTTP API (Medikations-Tagebuch)
# Startet nach dem Netzwerk, damit das iPad/Android zugreifen kann.
Wants=network-online.target
After=network-online.target

[Service]
# ── Verzeichnis, in das `npm run deploy` den Build installiert hat ──
WorkingDirectory=%h/mediary

# ── Start-Skript (erzeugt das Datenverzeichnis, falls nötig) ──
ExecStart=%h/mediary/start.sh

# ── Env ──
# DB und DEFAULTS liegen in ~/.local/share/mediary (server/src/config.ts)
# Override hier möglich, z. B.:
# Environment="PORT=4000"
# Environment="DB_PATH=/custom/path/mediary.db"
# Environment="DEFAULTS_PATH=/custom/path/DEFAULTS.md"

# ── Safety ──
Restart=on-failure
RestartSec=5s
# Kein Root – der User, der `systemctl --user` aufruft.
User=%u

# ── Logging ──
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mediary

[Install]
WantedBy=default.target