# mediary.service
# ─────────────────────────────────────────────────────────────────────────────
# meDiary systemd service (user service)
# Installieren:
#   cp mediary.service ~/.config/systemd/user/mediary.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now mediary
#
# Logs:  journalctl --user-unit=mediary -f
# Stop:  systemctl --user stop mediary
# Restart: systemctl --user restart mediary
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=meDiary HTTP API (Medikations-Tagebuch)
# Startet nach dem Netzwerk, damit das iPad/Android zugreifen kann.
Wants=network-online.target
After=network-online.target

[Service]
# ── Verzeichnis, in das `npm run deploy` den Build installiert hat ──
WorkingDirectory=%h/mediary

# ── Start-Skript (erzeugt das Datenverzeichnis, falls nötig) ──
ExecStart=%h/mediary/start.sh

# ── Env ──
# DB und DEFAULTS liegen in ~/.local/share/mediary (server/src/config.ts).
# Diese Env-Vars werden beim `npm run deploy` aus der .env im Projekt-Root
# hier injiziert. deploy.sh ersetzt die untenstehende Marker-Zeile durch
# die tatsächlichen `Environment=...`-Zeilen (eine pro Variable).
# Marker-Zeile (muss genau so aussehen, eine Zeile, eine Zeile, ein Match):
# __MEDIARY_INJECT_ENV_HERE__

# ── Safety ──
Restart=on-failure
RestartSec=5s
# Kein Root – der User, der `systemctl --user` aufruft.
User=%u

# ── Logging ──
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mediary

[Install]
WantedBy=default.target